# samples/alpine/nginx-php/7.4/etc/s6-overlay/s6-rc.d/crontab/finish
#!/bin/sh
# Cleanup script for crontab service
# Kill any remaining cron processes
pkill crond || true
