# samples/alpine/nginx-php/7.4/etc/s6-overlay/s6-rc.d/php-fpm/finish
#!/bin/sh
# Clean up PHP-FPM socket
rm -f /home/another/run/php-fpm.sock