# samples/alpine/nginx-php/7.4/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/bin/sh
exec s6-envuidgid another /usr/local/bin/php-fpm --nodaemonize --fpm-config /etc/php/php-fpm.conf