# samples/alpine/nginx-php/7.4/etc/s6-overlay/s6-rc.d/nginx/finish
#!/bin/sh
# Gracefully shutdown nginx as non-root user
nginx -s quit
