# samples/alpine/nginx-php/7.4/etc/s6-overlay/s6-rc.d/nginx/run
#!/bin/sh
# Use capabilities to allow nginx to bind to port 80 as non-root user
exec nginx -g "daemon off;"
